# build/crog.service
[Unit]
Description=crog
After=network.target

[Service]
ExecStart=/usr/bin/crog work --config /etc/crog/crog.yaml
StandardOutput=inherit
StandardError=inherit
Restart=always
User=root

[Install]
WantedBy=multi-user.target
