# build/healthcheck.service
[Unit]
Description=healthcheck
After=network.target

[Service]
ExecStart=/usr/bin/healthcheck work --config /etc/healthcheck/healthcheck.yaml
StandardOutput=inherit
StandardError=inherit
Restart=always
User=root

[Install]
WantedBy=multi-user.target
